# service/svscan/run
#!/bin/sh

/usr/sbin/dimq -c /etc/dimq/dimq.conf
